# systemd/mcp-gpt-expert.service.template
[Unit]
Description=MCP GPT Expert Server
After=network.target

[Service]
Type=simple
WorkingDirectory={{PROJECT_DIR}}
EnvironmentFile={{PROJECT_DIR}}/.env
ExecStart=/usr/bin/node {{PROJECT_DIR}}/dist/index.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
